# btl-server.service
[Unit]
Description=BTL Game Server
After=network.target

[Service]
ExecStart=/opt/btl/btl-server
WorkingDirectory=/opt/btl
Restart=on-failure
RestartSec=3

# Optional: override game port (cert hash HTTP served on game_port+1)
# Environment=BTL_PORT=5888

[Install]
WantedBy=multi-user.target
